# airootfs/etc/NetworkManager/dispatcher.d/99-connectivity-change.sh
#!/bin/bash
if [ "$2" = "connectivity-change" ]; then
    /usr/bin/fix-repository &
fi